# src/main/bin/conveyor-sets.ini
[program:conveyor-sets]
command=/opt/conveyor-sets/startup.sh
directory=/opt/conveyor-sets
user=root
autostart=true
autorestart=true
stopsinal=INT
stdout_logfile=/var/log/conveyor-sets.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
redirect_stderr=true
